# templates/backup_mastodon.service
[Unit]
Description=Service Mastodon Backup

[Service]
Type=simple
User=root
ExecStart=/opt/backup_mastodon.sh

[Install]
WantedBy=multi-user.target